# private_dot_config/systemd/user/goimapnotify@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/home/tetov/.local/bin/goimapnotify -conf %h/.config/imapnotify/%i.conf
