# private_dot_config/systemd/user/sync_emails@often.timer.d/OnUnitInactiveSec.conf
[Timer]
OnUnitInactiveSec=
OnUnitInactiveSec=2m
