# private_dot_config/systemd/user/git_sync@\x2elocal-share-chezmoi.timer.d/override.conf
[Timer]
OnUnitInactiveSec=
OnUnitInactiveSec=1d
