# private_dot_config/systemd/user/emacs_restart.service
[Unit]
Description=Restart emacs daemon

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart
