# private_dot_config/systemd/user/git-sync-org.timer
[Unit]
Description=git sync org dir periodically

[Timer]
OnStartupSec=5 m
OnUnitInactiveSec=10 m
Persistent=true

[Install]
WantedBy=timers.target
