# private_dot_config/systemd/user/git_sync@.timer
[Unit]
Description=git sync dir periodically

[Timer]
OnStartupSec=2m
OnUnitInactiveSec=10m
Persistent=true

[Install]
WantedBy=timers.target
