# private_dot_config/systemd/user/sync_emails@.timer
[Unit]
Description=Timer to synx mailboxes in group %i

[Timer]
OnBootSec=3m
OnCalendar=8m

[Install]
WantedBy=timers.target
