# private_dot_config/systemd/user/sync_emails@.timer
[Unit]
Description=Timer to synx mailboxes in group %i

[Timer]
OnBootSec=2m
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target
