# private_dot_config/systemd/user/sync_emails@.service
[Unit]
Description=Synchronise mailboxes in group %i

[Service]
Type=oneshot
ExecStart=%h/bin/sync_emails.sh %i
